# deploy/supervisord_docker.conf
[supervisord]
nodaemon=true

[program:aasaan]
command=/usr/local/bin/aasaan_start_docker

[program:aasaan-rq]
command=/usr/local/bin/aasaan_worker_start_docker
